# etc/systemd/system/screenshot.service.d/wayland.conf
# Wayland/labwc: grim compositor capture + MPV vo=gpu IPC (no kmsgrab/TTY).
# Install: /etc/systemd/system/screenshot.service.d/wayland.conf
# Requires: apt install grim

[Unit]
After=
Wants=
After=dsign-compositor.service dsign-mpv-wayland.service
Wants=dsign-compositor.service dsign-mpv-wayland.service

[Service]
EnvironmentFile=-/var/lib/dsign/config/wayland.env
Environment=DSIGN_DISPLAY_BACKEND=wayland
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000
# grim needs the compositor socket; kmsgrab DeviceAllow from base unit is unused.
DeviceAllow=
BindReadOnlyPaths=/run/user/1000
ReadWritePaths=/home/dsign/dsign/static/images/ /var/lib/dsign/mpv/ /var/tmp
