# etc/systemd/system/digital-signage.service.d/wayland.conf
# Enable with: sudo systemctl enable dsign-compositor dsign-logo dsign-mpv-wayland
# Disable DRM stack: sudo systemctl disable --now dsign-mpv.service
[Service]
Environment=DSIGN_DISPLAY_BACKEND=wayland
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Unit]
After=dsign-compositor.service dsign-mpv-wayland.service
Wants=dsign-compositor.service dsign-mpv-wayland.service dsign-logo.service
