# etc/systemd/system/digital-signage.service.d/wayland.conf
# Wayland stack: replaces DRM deps from digital-signage.service (Wants=dsign-mpv.service).
# Install to: /etc/systemd/system/digital-signage.service.d/wayland.conf
#   (note hyphen: digital-signage, NOT digital.signage)
#
# Enable:  sudo systemctl enable dsign-compositor dsign-logo dsign-mpv-wayland
# Disable DRM: sudo systemctl disable --now dsign-mpv.service && sudo systemctl mask dsign-mpv.service

[Service]
Environment=DSIGN_DISPLAY_BACKEND=wayland
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Unit]
# Clear base unit After=/Wants= so dsign-mpv.service is not pulled in.
After=
Wants=
After=network.target dsign-compositor.service dsign-mpv-wayland.service dsign-network-assistant.service
Wants=dsign-compositor.service dsign-mpv-wayland.service dsign-logo.service
