# etc/systemd/system/digital-signage.service
[Unit]
Description=Digital Signage Controller
After=network.target dsign-mpv.service
Requires=dsign-mpv.service

[Service]
Type=simple
User=dsign
Group=dsign
WorkingDirectory=/home/dsign
ExecStart=/usr/bin/python3 /home/dsign/run.py
# Настройки перезапуска
Restart=on-failure
RestartSec=30s
StartLimitIntervalSec=60
StartLimitBurst=3

# Настройки завершения
TimeoutStopSec=30
KillMode=mixed

# Ограничения ресурсов
MemoryLimit=500M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
