# etc/systemd/system/digital-signage.service
[Unit]
Description=Digital Signage Controller
After=network.target dsign-mpv.service
Wants=dsign-mpv.service

[Service]
User=dsign
Group=dsign
WorkingDirectory=/home/dsign
ExecStart=/usr/bin/python3 /home/dsign/run.py
Restart=on-failure
RestartSec=30s
KillMode=mixed
TimeoutStopSec=30

# Общие переменные среды
Environment=MPV_IPC_SOCKET=/tmp/mpv-socket
Environment=DISPLAY=:0
Environment=PYTHONPATH=/home/dsign

# Ограничения ресурсов
MemoryLimit=500M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
