# etc/systemd/system/digital-signage.service
[Unit]
Description=Digital Signage Controller
After=network.target dsign-mpv.service dsign-network-assistant.service
Requires=dsign-mpv.service
Wants=dsign-network-assistant.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=dsign
Group=dsign
# Default layout: git clone at /home/dsign/dsign → package /home/dsign/dsign/dsign/, run.py at /home/dsign/dsign/run.py
WorkingDirectory=/home/dsign/dsign
Environment=DSIGN_PROJECT_ROOT=/home/dsign/dsign
ExecStartPre=/bin/mkdir -p /var/lib/dsign/media/tmp /var/log/dsign
ExecStartPre=/bin/chown dsign:dsign /var/lib/dsign/media/tmp /var/log/dsign
ExecStartPre=/bin/chmod 775 /var/lib/dsign/media/tmp /var/log/dsign
# Boot race: Requires= starts us right after mpv's main PID is active; IPC socket may appear slightly later.
ExecStartPre=/bin/bash -c 'for _ in $(seq 1 120); do [[ -e /var/lib/dsign/mpv/socket ]] && exit 0; sleep 0.5; done; echo "MPV IPC socket missing after wait" >&2; exit 1'
Environment="TMPDIR=/var/lib/dsign/media/tmp"
Environment=DSIGN_LOG_LEVEL=WARNING
Environment=SOCKETIO_ASYNC_MODE=threading
Environment=DSIGN_AUTO_TRANSCODE_VIDEOS=true
Environment=DSIGN_TRANSCODE_TARGET_RESOLUTION=1920x1080
Environment=DSIGN_TRANSCODE_TARGET_FPS=25
# Optional: HLS/micro-hitch tuning for VK/Rutube (increases RAM/latency; can worsen jank on some Pi setups).
# Environment=DSIGN_MPV_NETBUF=1
# Environment=DSIGN_MPV_NETBUF_SECS=12
# Environment=DSIGN_MPV_NETBUF_MAX_MB=96
# Environment=DSIGN_MPV_NETBUF_BACK_MB=16
# Environment=DSIGN_MPV_NETBUF_READAHEAD_SECS=20
# Rare: re-enable strict IPC readiness wait for plain https streams (default off; Rutube HLS uses resolved URLs).
# Environment=DSIGN_MPV_STREAM_READY_HTTPS=1
# Seconds before treating idle-active as EOF for network streams (HLS can flicker idle while buffering; default 25).
# Environment=DSIGN_MPV_NETWORK_IDLE_GRACE_SEC=25
# Use project venv — system python has no Flask/deps (ModuleNotFoundError: flask).
ExecStart=/home/dsign/venv/bin/python /home/dsign/dsign/run.py
# Настройки перезапуска
Restart=on-failure
RestartSec=30s

# Настройки завершения
TimeoutStopSec=30
KillMode=mixed

# Ограничения ресурсов
# Загрузка в галерею до 1 GiB буферизуется в процессе; 500M обрывала крупные файлы (~70% от ~650 MB).
MemoryMax=2G
CPUQuota=50%

[Install]
WantedBy=multi-user.target
