# etc/systemd/system/dsign-network-assistant.service
[Unit]
Description=Digital Signage Network Assistant (screen mode)
After=network.target
Before=dsign-mpv.service digital-signage.service
Wants=network.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/bin/dsign-network-assistant
Environment=DSIGN_NETWORK_PROMPT_TIMEOUT_SEC=120
Environment=DSIGN_STARTUP_IP_FILE=/tmp/dsign-startup-ip.txt
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
