# etc/systemd/system/dsign-show-startup-ip.service
[Unit]
Description=Digital Signage Startup IP OSD helper
After=dsign-mpv.service
# Do not block digital-signage or mpv; this unit must finish in seconds.
Before=

[Service]
Type=oneshot
RemainAfterExit=no
TimeoutStartSec=10
User=dsign
Group=dsign
ExecStart=/usr/local/bin/dsign-show-startup-ip

[Install]
WantedBy=multi-user.target
