# etc/systemd/system/dsign-show-startup-ip.service
[Unit]
Description=Digital Signage Startup IP OSD helper
After=dsign-mpv.service dsign-mpv-wayland.service dsign-network-assistant.service

[Service]
Type=oneshot
RemainAfterExit=no
TimeoutStartSec=10
User=dsign
Group=dsign
Environment=DSIGN_STARTUP_IP_ALLOW_LIVE_LOOKUP=1
Environment=DSIGN_STARTUP_IP_SOCKET_WAIT_SEC=3
Environment=DSIGN_STARTUP_IP_BG_RETRY_SEC=60
ExecStart=-/bin/bash /usr/local/bin/dsign-show-startup-ip

[Install]
WantedBy=multi-user.target
