# etc/systemd/system/dsign-mpv.service
[Unit]
Description=MPV Player for Digital Signage
# Не указывать After=multi-user.target: оба юнита WantedBy=multi-user.target, а digital-signage
# Requires=dsign-mpv — получится цикл «multi-user → signage → mpv → после multi-user».
# MPV with DRM/KMS needs a controlling TTY/VT. Pi 3B+: --vo=drm is lighter than --vo=gpu.
After=network.target getty@tty1.service
Conflicts=getty@tty1.service

[Service]
User=dsign
Group=video
SupplementaryGroups=tty audio
PermissionsStartOnly=true
Environment="TERM=linux"
Environment="DSIGN_MPV_PROFILE=signboard"
WorkingDirectory=/home/dsign

StandardInput=tty
TTYPath=/dev/tty1
TTYReset=no
TTYVHangup=no
TTYVTDisallocate=no

# Root-only: chvt once per boot. '-' = ignore errors; bash fallback avoids 203/EXEC if script missing.
ExecStartPre=-/bin/bash -c 'if [ -x /usr/local/bin/dsign-mpv-pre-display ]; then exec /usr/local/bin/dsign-mpv-pre-display; else mkdir -p /run/dsign 2>/dev/null; /bin/chvt 1 2>/dev/null || true; fi'
ExecStartPre=/bin/mkdir -p /var/lib/dsign/mpv /var/lib/dsign/mpv/archive
ExecStartPre=/bin/chown dsign:video /var/lib/dsign/mpv /var/lib/dsign/mpv/archive
ExecStartPre=/bin/chmod 775 /var/lib/dsign/mpv
ExecStartPre=/bin/rm -f /var/lib/dsign/mpv/socket

ExecStart=/usr/local/bin/dsign-mpv-launch
ExecStartPost=-/bin/systemctl --no-block start dsign-show-startup-ip.service
ExecStopPost=-/bin/bash -c "tr -d '\\r' < /usr/local/bin/dsign-mpv-archive-log | bash"

Restart=always
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
