# etc/systemd/system/dsign-network-assistant.service
[Unit]
Description=Digital Signage Network Assistant (screen mode)
After=network.target
# Run early but do not block mpv/signage start (parallel boot).
Before=
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root
ExecStart=/usr/local/bin/dsign-network-assistant
# Optional per-device config (written by Settings UI). Must be safe to ignore.
EnvironmentFile=-/var/lib/dsign/config/network-assistant.env
Environment=DSIGN_NETWORK_PROMPT_TIMEOUT_SEC=120
Environment=DSIGN_STARTUP_IP_FILE=/tmp/dsign-startup-ip.txt

[Install]
WantedBy=multi-user.target
