# etc/systemd/system/dsign-mpv.service
[Unit]
Description=MPV Player for Digital Signage
After=graphical.target
Requires=graphical.target

[Service]
User=dsign
Group=video
Environment="DISPLAY=:0"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="PULSE_RUNTIME_PATH=/run/user/1000/pulse"
WorkingDirectory=/home/dsign

# Подготовка директории и прав
ExecStartPre=/bin/mkdir -p /var/lib/dsign/mpv
ExecStartPre=/bin/chown dsign:video /var/lib/dsign/mpv
ExecStartPre=/bin/chmod 775 /var/lib/dsign/mpv
ExecStartPre=/bin/rm -f /var/lib/dsign/mpv/socket

# Основная команда запуска (ИСПРАВЛЕН ПУТЬ)
ExecStart=/usr/bin/mpv --idle=yes --no-config --no-osc --no-input-default-bindings --input-ipc-server=/var/lib/dsign/mpv/socket --hwdec=auto --vo=drm --log-file=/var/log/dsign/mpv.log

# Настройки перезапуска
Restart=always
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
