# etc/systemd/system/dsign-wifi-on-display.service
[Unit]
Description=Digital Signage Wi-Fi picker on content display (release MPV DRM, nmtui)
After=dsign-mpv.service
Wants=dsign-mpv.service

[Service]
Type=oneshot
RemainAfterExit=no
User=root
Group=root
EnvironmentFile=-/var/lib/dsign/config/network-assistant.env
Environment=DSIGN_MPV_SOCKET=/var/lib/dsign/mpv/socket
Environment=DSIGN_NETWORK_STATUS_FILE=/run/dsign/network-status.env
Environment=DSIGN_NETWORK_WIFI_PROMPT_SEC=60
Environment=DSIGN_MPV_RESTORE_VO=drm
ExecStart=/usr/local/bin/dsign-wifi-on-display

[Install]
WantedBy=multi-user.target
