# etc/systemd/system/dsign-mpv-wayland.service
[Unit]
Description=MPV Player for Digital Signage (Wayland / vo=gpu)
Documentation=file:/usr/local/bin/dsign-mpv-launch-wayland
After=dsign-compositor.service dsign-logo.service network.target
Requires=dsign-compositor.service
Wants=dsign-logo.service dsign-network-assistant.service

[Service]
Type=simple
User=dsign
Group=video
SupplementaryGroups=video audio render
EnvironmentFile=-/var/lib/dsign/config/wayland.env
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DSIGN_DISPLAY_BACKEND=wayland
Environment=DSIGN_MPV_PROFILE=intel-balanced-wayland
ExecStartPre=/bin/mkdir -p /var/lib/dsign/mpv /var/lib/dsign/mpv/archive
ExecStartPre=/bin/chown dsign:video /var/lib/dsign/mpv /var/lib/dsign/mpv/archive
ExecStartPre=/bin/chmod 775 /var/lib/dsign/mpv
ExecStartPre=/bin/rm -f /var/lib/dsign/mpv/socket
ExecStart=/usr/local/bin/dsign-mpv-launch-wayland
ExecStartPost=-/bin/systemctl --no-block restart dsign-show-startup-ip.service
ExecStartPost=-/bin/bash /usr/local/bin/dsign-mpv-post-start
ExecStopPost=-/bin/bash -c "tr -d '\\r' < /usr/local/bin/dsign-mpv-archive-log | bash"
Restart=always
RestartSec=5s
StartLimitInterval=300s
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
