# etc/systemd/system/screenshot.service
[Unit]
Description=DSign Screenshot Service
After=dsign-mpv.service
Wants=dsign-mpv.service

[Service]
# Wait until dsign-capture exits so `systemctl start` means "capture finished"
# (Type=simple returns immediately and the UI can reload before the file exists).
Type=oneshot
User=dsign
Group=video
WorkingDirectory=/home/dsign/dsign/static/images
ExecStart=/usr/local/bin/dsign-capture
Restart=on-failure
RestartSec=5s
TimeoutStopSec=5
SuccessExitStatus=0 1

# Security sandboxing
# NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/dsign/dsign/static/images/

[Install]
WantedBy=multi-user.target
