# etc/systemd/system/dsign-mpv.service.d/99-intel.conf
[Service]
Environment=DSIGN_MPV_PROFILE=intel-balanced
