# etc/systemd/system/dsign-mpv.service.d/intel-iris-xe.conf.example
# Intel Iris Xe / x86 — выбор профиля через переменную (mpv.conf не трогаем):
#
#   sudo mkdir -p /etc/systemd/system/dsign-mpv.service.d
#   sudo cp intel-iris-xe.conf.example /etc/systemd/system/dsign-mpv.service.d/intel-iris-xe.conf
#   sudo systemctl daemon-reload && sudo systemctl restart dsign-mpv.service
#
# По умолчанию в примере — баланс качество/нагрузка. Для максимума качества замените на:
#   Environment=DSIGN_MPV_PROFILE=intel-high-quality

[Service]
Environment=DSIGN_MPV_PROFILE=intel-balanced
